# scripts/windows/cryostat-wsl-boot-task.xml
<?xml version="1.0" encoding="UTF-8"?>
<!--
  cryostat-wsl-boot — boots the Ubuntu-24.04 WSL distro at Windows startup and
  keeps the VM alive (deployment-wsl.md §5). The long-lived `sleep infinity`
  pins the WSL2 VM up; systemd inside the distro then starts cryo-ingest,
  cryo-watchdog, and the backup timer.

  REPLACE `YOURPC\youruser` below with the Windows account that installed
  Ubuntu (the distro is registered under that profile - SYSTEM can't see it),
  or override it at import time with /ru. Import and verification commands are
  in README.md next to this file.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Boot the cryostat-monitor WSL2 distro (Ubuntu-24.04) at startup and keep the VM alive. See docs/deployment-wsl.md §5.</Description>
    <URI>\cryostat-wsl-boot</URI>
  </RegistrationInfo>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
      <Delay>PT30S</Delay>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <!-- PLACEHOLDER: the account that installed Ubuntu-24.04 -->
      <UserId>YOURPC\youruser</UserId>
      <!-- Password = "Run whether user is logged on or not" (stored credentials) -->
      <LogonType>Password</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <!-- PT0S = no execution time limit; the action runs forever by design. -->
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Windows\System32\wsl.exe</Command>
      <Arguments>-d Ubuntu-24.04 -u root -- sleep infinity</Arguments>
    </Exec>
  </Actions>
</Task>
